# config/systemd/redash-worker.service
[Unit]
Description=Redash Background Worker
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/redash
ExecStart=/opt/redash/bin/worker
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
